# s6-rc.d/metrics_exporter/run
#!/command/execlineb -P
with-contenv
/app/wg_metrics_exporter
