# s6-rc.d/node_service_log_prepare/up
if { mkdir -p /var/log/node_service }
if { chown nobody:nogroup /var/log/node_service }
chmod 02755 /var/log/node_service
