# s6-rc.d/node_service/run
#!/command/execlineb -P
with-contenv
exec node /app/server.js 2>&1
